# infra/roles/playground/templates/playground.service.j2
[Unit]
Description=Log Collector Playground (docker compose stack)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory={{ app_dir }}
EnvironmentFile={{ app_dir }}/.env
# Don't `docker compose pull` — the `app` service is built locally and has no
# registry image. `up --build` handles both building `app` and pulling `caddy`.
ExecStart=/usr/bin/docker compose up -d --build --remove-orphans
ExecStop=/usr/bin/docker compose down
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
